# etc/cron.d/boxy.umcg-gaf-dm.cron
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; copyRawDataToPrm.sh     -g umcg-gaf -s gattaca02.gcc.rug.nl"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; copyProjectDataToPrm.sh -g umcg-gaf"
