# etc/cron.d/chaperone.umcg-gd-dm.cron
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.12-bare; copyRawDataToPrm.sh -g umcg-gd -s gattaca01.gcc.rug.nl -l WARN"
