# etc/cron.d/gattaca01.umcg-gd-ateambot.cron
0 6 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; ${HOME}/FileSystemCheckGD.sh"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load ngs-utils/19.03.3 ; checkSamplesheet.sh -s /groups/umcg-gd/scr01/"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-NGS_Demultiplexing-2.3.1; demultiplexing.sh -g umcg-gd -l WARN"
*/15 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; notifications.sh -g umcg-gd -e"
*/15 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; demultiplexingTiming.sh -g umcg-gd"
