# etc/cron.d/gattaca01.umcg-atd-ateambot.cron
*/1 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.12-NGS_Demultiplexing-2.2.12; demultiplexing.sh -g umcg-atd"
