# etc/cron.d/leucine-zipper.umcg-gd-ateambot.cron
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-NGS_DNA-3.5.5; startPipeline.sh -g umcg-gd"
*/5 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; notifications.sh -g umcg-gd -e"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-NGS_DNA-3.5.5; calculateProjectMd5s.sh  -g umcg-gd"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_DNA/3.5.5; GavinSA.sh -g umcg-gd -t tmp06"
0,15,30,45 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; pipelineTiming.sh -g umcg-gd -s gattaca01.gcc.rug.nl -r /groups/umcg-gd/scr01/"
10,25,40,55 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; pipelineTiming.sh -g umcg-gd -s gattaca02.gcc.rug.nl -r /groups/umcg-gd/scr01/"
0 1 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; ConcordanceMakeSamplesheet.sh -g umcg-gd -a umcg-gap"
0 2 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; ConcordanceCheck.sh -g umcg-gd"
0 6 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; ${HOME}/FileSystemCheckGD.sh"
