# etc/cron.d/leucine-zipper.umcg-gd-ateambot.cron
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.12-NGS_DNA-3.4.4; startPipeline.sh -g umcg-gd"
*/5 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.12-bare; notifications.sh -g umcg-gd -e"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; sh /groups/umcg-gd/tmp06/GavinStandAlone/cleanup.sh umcg-gd"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; sh /groups/umcg-gd/tmp06/GavinStandAlone/GavinSA.sh"
0 6 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; ${HOME}/FileSystemCheckGD.sh"
*/1 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; bash /groups/umcg-gd/tmp06/logs/correctPipeline.sh > /tmp/test4.txt 2>&1"
