# etc/cron.d/leucine-zipper.umcg-gap-ateambot.cron
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-GAP-2.3.1; startPipeline.sh -g umcg-gap"
*/5 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-bare; notifications.sh -g umcg-gap -e"
*/10 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/3.0.1-GAP-2.3.1; calculateProjectMd5s.sh  -g umcg-gap"
