# etc/cron.d/zinc-finger.umcg-genomescan-ateambot.cron
*/30 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.10-bare; processGsRawData.sh     -g umcg-genomescan"
*/15 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load NGS_Automated/2.0.10-bare; notifications.sh        -g umcg-genomescan -e"